# bootstrap-systems/spoofie-kvm-vm9333524/rootfs-changes/etc/rc.local
#!/bin/bash
#
# /etc/rc.local: local multi-user startup script
#

localedef -i "en_US" -f "UTF-8" "en_US.UTF-8"

echo -e "toor\ntoor" | (passwd root)

# End of file
